# hass-cert-update.cron
*/5 * * * * root /path/to/hass-cert-update --host your.domain.com 